# config/supervisor/mongo.conf
[program:mongo]
command=mongod
user=root
directory=/data/db
environment=
autostart=true
autorestart=false
