# config/supervisor/mysql.conf
[program:mysql]
command=/usr/local/bin/mysql_start.sh
user=root
directory=/root
environment=MYSQL_DATABASE=openedx,MYSQL_USER=openedx,MYSQL_PASSWORD=password
autostart=true
autorestart=false
